# .config/systemd/user/atuin.service
[Unit]
Description=Atuin Daemon

[Service]
Restart=always
ExecStart=/home/%u/.local/share/cargo/bin/atuin daemon

[Install]
WantedBy=default.target
